# src/autostart/linux.desktop
[Desktop Entry]
Type=Application
Exec={{EXEC}}
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=MyAppStart
Comment=Starts MyApp
